# fs/etc/services.d/connector/run
#!/usr/bin/with-contenv sh

exec runas /swarm-local-network-connector